# systemd/homelab-backup.timer
[Unit]
Description=Daily Home Lab Backup Timer
Requires=homelab-backup.service

[Timer]
OnCalendar=daily
# Run at 2:00 AM every day
OnCalendar=*-*-* 02:00:00
# If the system was off, run immediately when it comes back on
Persistent=true
# Randomize start time by 0-30 minutes to avoid load spikes
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
