# systemd/homelab-backup.service
[Unit]
Description=Home Lab Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root
# IMPORTANT: Update this path to the actual location of backup.sh
# Example: ExecStart=/usr/bin/bash /home/pi/homelab/backup.sh
ExecStart=/usr/bin/bash /path/to/homelab/backup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
